# systemd/sysc-walls.service
[Unit]
Description=sysc-walls idle daemon and screensaver
After=graphical-session.target
Documentation=https://github.com/Nomadcxx/sysc-walls

[Service]
Type=simple
ExecStart=/usr/local/bin/sysc-walls-daemon -start -debug
Restart=on-failure
RestartSec=5
TimeoutStartSec=10
TimeoutStopSec=10

# Environment variables for user session
Environment=DISPLAY=%E{DISPLAY}
Environment=WAYLAND_DISPLAY=%E{WAYLAND_DISPLAY}
Environment=XDG_RUNTIME_DIR=%t
Environment=HYPRLAND_INSTANCE_SIGNATURE=%E{HYPRLAND_INSTANCE_SIGNATURE}

[Install]
WantedBy=graphical-session.target